# packaging/systemd/chitra-triaged.service.example
[Unit]
Description=chitra triaged (state-transition dedup daemon)
Documentation=https://github.com/first-polyphony/chitra
After=network-online.target

[Service]
Type=simple
# Replace with the user/group this service should run as. Do not run as root.
User=chitra
Group=chitra
# systemd creates and owns this directory (mode 0750) for the unit's lifetime.
StateDirectory=chitra
# Fill in the venv path where `pip install chitra` (or `pip install -e .`) was run.
ExecStart=/path/to/venv/bin/triaged --events-log /var/lib/chitra/events.log --state-file /var/lib/chitra/triaged-state.json --triage-log /var/lib/chitra/triaged.log
Restart=always
RestartSec=10
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
